# rc.netisdn.in
#! /bin/sh
##
##    A NetISDN server init script for LSB-compliant Linux systems.
##
##    (C) Copyright 2008 Fred Gleason <fredg@paravelsystems.com>
##
##      $Id: rc.netisdn.in,v 1.1 2008/09/05 18:38:10 cvs Exp $
##
##    This program is distributed in the hope that it will be useful,
##    but WITHOUT ANY WARRANTY; without even the implied warranty of
##    MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
##    GNU General Public License for more details.
##
##    You should have received a copy of the GNU General Public License
##    along with this program; if not, write to the Free Software
##    Foundation, Inc., 59 Temple Place, Suite 330, 
##    Boston, MA  02111-1307  USA
##

### BEGIN INIT INFO
# Provides:          netisdn
# Required-Start:    $remote_fs mysql
# Required-Stop:     $remote_fs
# Default-Start:     2 3 5
# Default-Stop:      0 1 6
# Description:       Start the NetISDN server system daemons
### END INIT INFO


function StartDaemons {
  @LOCAL_PREFIX@/sbin/netsipd
}
function StopDaemons {
  kill `cat /var/run/netsipd.pid` > /dev/null 2> /dev/null
  rm -f /var/run/netsipd.pid
}

# Set path for script functions
export PATH=$PATH:/usr/local/bin:/usr/local/sbin

#
# If the LSB Init functions exist, use them, otherwise roll our own
#
if test -f /lib/lsb/init-functions ; then
  . /lib/lsb/init-functions
else
  function log_success_msg() {
    echo "$1 ... done."
  }
  function log_failure_msg() {
    echo "$1 ... done."
  }
fi

# Check for missing binaries
if [ ! -x @LOCAL_PREFIX@/sbin/netsipd ] ; then
  echo "netsipd not installed"
  exit 5
fi

# Check for existence of needed config file and read it
if [ ! -r /etc/netisdn.conf ] ; then
  echo "Missing /etc/netisdn.conf"
  exit 6
fi


case "$1" in
    start)
	$0 silent-status
	if test $? -eq 3 ; then
	    StartDaemons
	fi
	$0 silent-status
	if test $? -eq 0 ; then
	    log_success_msg "Starting NetISDN system daemon"
	    exit 0
	else
	    log_failure_msg "Starting NetISDN system daemon"
	    exit 1
	fi
	;;
    stop)
	StopDaemons
	$0 silent-status
	if test $? -eq 0 ; then
	    log_failure_msg "Stopping NetISDN system daemon"
	    exit 1
	else
	    log_success_msg "Stopping NetISDN system daemon"
	    exit 0
	fi
	;;
    restart)
	$0 silent-status
	if test $? -eq 0 ; then
	    StopDaemons
	fi
	$0 silent-status
	if test $? -eq 0 ; then
	    exit 1
	fi
	StartDaemons
	$0 silent-status
	if test $? -eq 0 ; then
	    exit 0
	fi
	exit 1
	;;
    force-reload)
	$0 restart
	;;
    status)
	$0 silent-status
	EXIT_CODE=$?
	if test $EXIT_CODE -eq 0 ; then
	    echo "NetISDN system daemon running."
	else 
	    echo "NetISDN system daemon stopped."
	fi
	exit $EXIT_CODE
	;;
    silent-status)
	if test -f /var/run/netsipd.pid ; then
	    if test -d /proc/`cat /var/run/netsipd.pid` ; then
		exit 0
	    fi
	fi
	exit 3
	;;
    *)
	echo "Usage: $0 {start|stop|restart|force-reload|status}"
	exit 1
	;;
esac
