# systemd/netisdn.service.in
[Unit]
Description=NetISDN Presence Server
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=simple
ExecStart=@prefix@/sbin/netsipd
PrivateTmp=true
Restart=always
RestartSec=2
StartLimitInterval=120
StartLimitBurst=50

[Install]
WantedBy=multi-user.target
